# ops/com.ai-job-tracker.ollama.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai-job-tracker.ollama</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/ollama</string>
    <string>serve</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OLLAMA_HOST</key>
    <string>127.0.0.1:11434</string>
    <key>OLLAMA_CONTEXT_LENGTH</key>
    <string>4096</string>
    <key>OLLAMA_NUM_PARALLEL</key>
    <string>2</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/ai-job-tracker-ollama.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/ai-job-tracker-ollama.err</string>
</dict>
</plist>
